# deploy/hermes-mariyam.service
# Hermes/Mariyam backend — systemd unit (TEMPLATE, не использовать как есть)
#
# Установка на VPS (Ubuntu 24.04), только после копирования реальных секретов:
#   sudo install -d -m 700 /opt/hermes-mariyam-secrets
#   sudo install -m 600 /path/to/real-backend.env /opt/hermes-mariyam-secrets/backend.env
#   sudo cp deploy/hermes-mariyam.service /etc/systemd/system/
#   sudo systemd-analyze verify /etc/systemd/system/hermes-mariyam.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hermes-mariyam.service
#
# Реальные токены/пароли НЕ вписывать сюда; они в
# /opt/hermes-mariyam-secrets/backend.env. EnvironmentFile кладёт их в env процесса,
# docker compose использует интерполяцию переменных окружения.
#
# Type=oneshot + RemainAfterExit=yes: systemd только поднимает/останавливает compose.
# Перезапуск контейнеров обеспечивает docker compose restart: unless-stopped + docker.service.
# НЕ трогать /opt/time-agent, time_agent_bot, Time-Agent .env, SQLite volume, logs, backups.

[Unit]
Description=Hermes/Mariyam PostgreSQL via docker compose
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hermes-mariyam
EnvironmentFile=/opt/hermes-mariyam-secrets/backend.env
ExecStart=/usr/bin/docker compose up -d hermes_mariyam_postgres
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
